# templates/default/neo4j.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.neo4j.server</string>
	<key>WorkingDirectory</key>
  <string><%= @neo4j_prefix -%>/libexec</string>
	<key>ProgramArguments</key>
	<array>
  <string><%= @neo4j_prefix -%>/libexec/bin/neo4j</string>
		<string>console</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>OnDemand</key>
	<true/>
</dict>
</plist>
